# systemd/user/dms.service.d/override.conf
[Service]
Environment="PATH=/run/current-system/sw/bin:/usr/bin:/bin"
Environment="QT_QPA_PLATFORM=wayland"
Environment="WAYLAND_DISPLAY=wayland-1"
Environment="XDG_SESSION_TYPE=wayland"